# autostart/idle-check.service
[Unit]
Description=Check for idle session and start recording

[Service]
ExecStart=/home/pi/autostart/idle-check.sh
